# server/joystick.service
[Unit]
Description=Marcel Joystick

[Service]
Type=simple
ExecStart=/usr/bin/marcel-joystick

[Install]
WantedBy=multi-user.target
